# csr/crontab
# only on CSR_BASE
0 3 * * * "$HOME/.crontab-daily" > "$HOME/.crontab-daily-log" 2>&1
